# packaging/systemd/grimoired.service.in
[Unit]
Description=Grimoire voice command daemon
Documentation=https://github.com/salvogendut/grimoire
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
EnvironmentFile=-%h/.config/grimoire/grimoired.env
ExecStart=@BINDIR@/grimoired --listen-service --record-seconds 3
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
